# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:pyload]
command=python /opt/pyload/pyLoadCore.py --configdir=/opt/.pyload --daemon